# systemd/amplifi-backend.service
[Unit]
Description=ampliFi Cleanup Backend (FastAPI)
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=macan
Group=macan
WorkingDirectory=/home/macan/Development/ampliFi-center-management/backend
EnvironmentFile=/home/macan/Development/ampliFi-center-management/.env
ExecStart=/home/macan/Development/ampliFi-center-management/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port ${BACKEND_PORT}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
